# deploy/systemd/hadr-monitor.service
[Unit]
Description=HADR monitor — feed poller (USGS + GDACS + ReliefWeb)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust to where you cloned the repo and which user runs it.
WorkingDirectory=/opt/hadr-monitor
User=hadr
ExecStart=/usr/bin/env uv run hadr run
Restart=on-failure
RestartSec=5
# Optional: keep secrets/config out of the unit.
# EnvironmentFile=/opt/hadr-monitor/.env

[Install]
WantedBy=multi-user.target
